# flatpak.timer
[Unit]
Description=Run flatpak update weekly and on boot

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
